# deploy/llm-router.service
[Unit]
Description=LLM Router (ds4 <-> Ollama switchboard on :8077)
# Installed as a *user* unit (~/.config/systemd/user/) so it shares the same
# `systemctl --user` manager as ds4.service and can start/stop it. Lingering
# (loginctl enable-linger) makes it start at boot without a login session.
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/grahamfm/llm-router
Environment=PYTHONUNBUFFERED=1
Environment=ROUTER_CONFIG=/home/grahamfm/llm-router/config.yaml
ExecStart=/usr/bin/python3 -m router --config /home/grahamfm/llm-router/config.yaml
Restart=on-failure
RestartSec=3
TimeoutStopSec=15

[Install]
WantedBy=default.target
